# cookbooks/logstash_server/templates/default/logstash-upstart.conf.erb
#!upstart
description "logstash-server"

#env USER=<%= node[:logstash][:user] %>

env USER=logstash

start on startup

stop on shutdown

respawn

exec sudo -u $USER bash -c "/opt/logstash-1.4.2/bin/logstash --verbose -f /etc/logstash-1.4.2/server.conf \
   2>&1 >> /var/log/logstash.log"

# <%= node[:logstash][:log_dir] %>/logstash.log"